# nemo-webex-bot.service
[Unit]
Description=Nemo Webex Bot - LLM + MCP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=penta
WorkingDirectory=/home/penta/nemo-webex-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/penta/nemo-webex-bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
